# .crontab
SHELL=/bin/bash
0 7 * * * /usr/bin/minicpan >/tmp/minicpan.log 2>&1

