# .crontab
SHELL=/bin/bash
0 4 * * * /home/mtk/.bin/cs-update.sh

